# root/etc/s6-overlay/s6-rc.d/bot-service/run
#!/usr/bin/with-contenv bash
# permissions
mkdir -p /config
lsiown -R abc:abc /config /app
chmod -R 775 /config /app

exec s6-setuidgid abc python3 /app/app/main.py
